# webhook-server.service
[Unit]
Description=Personal Website Webhook Server
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=igor
ExecStart=/usr/bin/node /home/igor/personal-site/webhook-server.js
WorkingDirectory=/home/igor/personal-site
Environment=NODE_ENV=production
Environment=WEBHOOK_PORT=3001
Environment=WEBHOOK_SECRET=your-webhook-secret-change-this

# Logging
StandardOutput=append:/var/log/webhook-server.log
StandardError=append:/var/log/webhook-server.log

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/igor/personal-site /var/log
PrivateTmp=yes

[Install]
WantedBy=multi-user.target